# flask_blog/blog.conf
[program:blog]
command=/root/.virtualenvs/env_test/bin/gunicorn wsgi --bind 0.0.0.0:5000 --pid /tmp/blog.pid
directory=/root/Flask_Blog_Linux/flask_blog
autostart=true
